# cron/browseros-bookmark-refresh/com.max.browseros-bookmark-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.max.browseros-bookmark-refresh</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/max/.agents/cron/browseros-bookmark-refresh/browseros-bookmark-refresh.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/max/Library/Logs/browseros-bookmark-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/max/Library/Logs/browseros-bookmark-refresh.err.log</string>
  </dict>
</plist>
