# cron/nightshift/com.max.nightshift.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.max.nightshift</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/max/.agents/cron/nightshift/nightshift.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1200</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/max</string>
    <key>LANG</key>
    <string>en_US.UTF-8</string>
    <key>PATH</key>
    <string>/Users/max/Library/pnpm:/Users/max/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/max/.cron-logs/nightshift.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/max/.cron-logs/nightshift.log</string>
</dict>
</plist>
