# cron/cleanup-processes/com.max.cleanup-processes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.max.cleanup-processes</string>
	<key>Program</key>
	<string>/Users/max/.agents/cron/cleanup-processes/cleanup-processes.sh</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/max/Library/Logs/cleanup-processes.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/max/Library/Logs/cleanup-processes.log</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
